# config/supervisor/hydro-tcp-server.ini
[program:hydro-tcp-server]
command=php $PWD/bin/console tcp-server:start
startsecs=60
autostart=yes
autorestart=yes
