# config/supervisor/hydro-tcp-server.ini
[program:hydro-tcp-server]
command=APP_ENV=prod APP_DEBUG=0 php $PWD/bin/console tcp-server:start
startsecs=60
autostart=yes
autorestart=yes
